# meshtastic-dashboard.service
# Example unit for a host laid out under /opt/meshyface.
# scripts/deploy_meshyface.sh renders a host-specific unit during --bootstrap.
[Unit]
Description=Meshtastic Dashboard Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meshyface
Group=dialout
WorkingDirectory=/opt/meshyface
EnvironmentFile=/opt/meshyface/config/dashboard.env
ExecStart=/opt/meshyface/.venv/bin/python /opt/meshyface/app/mesh_dashboard.py --mesh-host ${MESH_HOST} --mesh-tcp-port ${MESH_PORT} --http-host ${DASH_HOST} --http-port ${DASH_PORT} --refresh-ms ${REFRESH_MS}
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
